# asus_touchpad.service
[Unit]
Description=Asus Touchpad to Numpad Handler

[Service]
Type=simple
ExecStart=/usr/share/asus_touchpad_numpad-driver/asus_touchpad.py $LAYOUT $PERCENTAGE_KEY > /var/log/asus_touchpad_numpad-driver/error.log 2>&1
StandardInput=tty-force
TimeoutSec=5

[Install]
WantedBy=default.target
